# wordpress/.docker/php/supervisor.d/conf.d/crond.conf
[program:crond]
process_name=%(program_name)s
command=/usr/sbin/crond -f -l 8
autostart=true
autorestart=true
user=docker
redirect_stderr=true
stdout_logfile=/var/log/cron/cron.log
